# deploy/webmenu.service
[Unit]
Description=App Menu (Next.js)
After=network.target

[Service]
Type=simple
User=<username>
WorkingDirectory=/home/<username>/webmenu
Environment=NODE_ENV=production
Environment=NVM_DIR=/home/<username>/.nvm
ExecStart=/bin/bash -lc "source /home/<username>/.nvm/nvm.sh && exec npm start -- -p 4000"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
